# scripts/ans-dashboard.service
[Unit]
Description=ANS Dashboard client (Vite) + API server
After=network.target
StartLimitIntervalSec=500

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/projetos/indicadores/ans-dashboard
Environment=NODE_ENV=production
Environment=SERVER_PORT=4000
ExecStart=/usr/bin/env bash /home/ubuntu/projetos/indicadores/ans-dashboard/scripts/start-dashboard.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
